# .config/systemd/user/plasma-powerdevil.service.d/override.conf
[Service]
Environment=POWERDEVIL_NO_DDCUTIL=1
